# recipes-pomelo-iot/pomelo-mqtt-dev/files/pomelo-mqtt-dev.service
[Unit]
Description=Pomelo MQTT server

[Service]
ExecStart=/usr/bin/hello

[Install]
WantedBy=multi-user.target
